# configs/fim-custom-rules.conf
<!--
  local_rules.xml — Custom Wazuh Detection Rules
  Author: Kshitiz
  Location: /var/ossec/etc/rules/local_rules.xml
  
  Custom rule IDs must be between 100000-119999 to avoid
  conflicts with Wazuh built-in rules.
-->

<group name="local,custom,">

  <!-- ============================================================
       RULE: Detect file creation in /tmp
       Parent: Rule 554 (new file added to system)
       Level 7 = medium severity
  ============================================================ -->
  <rule id="100001" level="7">
    <if_sid>554</if_sid>
    <field name="file">/tmp/</field>
    <description>Suspicious: New file created in /tmp directory</description>
    <group>file_creation,suspicious,tmp,</group>
    <mitre>
      <id>T1059</id>
    </mitre>
  </rule>

  <!-- ============================================================
       RULE: Detect changes to sudoers file
       Level 12 = high severity — sudoers changes can grant root
  ============================================================ -->
  <rule id="100002" level="12">
    <if_sid>550,554</if_sid>
    <field name="file">/etc/sudoers</field>
    <description>CRITICAL: sudoers file has been modified</description>
    <group>syscheck,sudoers,privilege_escalation,</group>
    <mitre>
      <id>T1548.003</id>
    </mitre>
  </rule>

  <!-- ============================================================
       RULE: Detect new SSH authorized_keys file
       An attacker adding their key achieves persistence
  ============================================================ -->
  <rule id="100003" level="10">
    <if_sid>554</if_sid>
    <field name="file">authorized_keys</field>
    <description>New SSH authorized_keys file detected — possible backdoor</description>
    <group>ssh,persistence,</group>
    <mitre>
      <id>T1098.004</id>
    </mitre>
  </rule>

  <!-- ============================================================
       RULE: Detect crontab modifications
       Attackers modify crontabs for persistence
  ============================================================ -->
  <rule id="100004" level="9">
    <if_sid>550,554</if_sid>
    <field name="file">/etc/cron</field>
    <description>Crontab file modified — check for persistence mechanisms</description>
    <group>cron,persistence,</group>
    <mitre>
      <id>T1053.003</id>
    </mitre>
  </rule>

  <!-- ============================================================
       RULE: Detect /etc/passwd modification
       Adding a user entry directly to /etc/passwd is suspicious
  ============================================================ -->
  <rule id="100005" level="11">
    <if_sid>550</if_sid>
    <field name="file">/etc/passwd</field>
    <description>ALERT: /etc/passwd has been directly modified</description>
    <group>syscheck,passwd,account_manipulation,</group>
    <mitre>
      <id>T1136.001</id>
    </mitre>
  </rule>

  <!-- ============================================================
       RULE: Multiple FIM alerts in short time = possible attack
       Parent: Rule 550 (file modified)
       Frequency: 10 times in 120 seconds
  ============================================================ -->
  <rule id="100006" level="12" frequency="10" timeframe="120">
    <if_matched_sid>550</if_matched_sid>
    <same_source_ip />
    <description>Multiple file integrity violations in short timeframe — possible attack</description>
    <group>syscheck,attack,</group>
  </rule>

</group>
